# resources/install/debian/postinst
#!/bin/bash
# postinst script for jitsi-videobridge

set -e

case "$1" in
    configure)

        CONFIG="/etc/jitsi/videobridge/config"

        # we don't want to regenerate config on upgrade
        OLDCONFIG="false"
        # migrate any old config found
        if [ -f "/etc/default/jitsi-videobridge" ]; then
            mv /etc/default/jitsi-videobridge $CONFIG
        fi

        if [ -f $CONFIG ]; then
            . $CONFIG
            if [ -n "$JVB_HOSTNAME" ] && [ -n "$JVB_PORT" ] && [ -n "$JVB_SECRET" ]; then
                OLDCONFIG="true"
            fi
        fi

        # generate config on new install
        if [ "$OLDCONFIG" = "false" ]; then

            # debconf hostname question
            . /usr/share/debconf/confmodule
            LOCALNAME="$(hostname)"
            LOCALIP="$(hostname -I)"
            NAMEBYIP="$(host $LOCALNAME | cut -f4 -d ' ')"

            # FIXME -- we have to check FQDN
            NAMEVALID=false
            for IP in $LOCALIP; do
                if [ "$NAMEBYIP" = "$IP" ]; then
                    NAMEVALID=true
                fi
            done

            # the default vars
            db_get jitsi-videobridge/jvb-hostname
            JVB_HOSTNAME="$RET"
            # 8-chars random secret, alternative to pwgen 8
            JVB_SECRET=`head -c 8 /dev/urandom | tr '\0-\377' 'a-zA-Z0-9a-zA-Z0-9a-zA-Z0-9a-zA-Z0-9@@@@####'`
            JVB_PORT=5347

            # storing default for later use by Jitsi Meet and other packages
            echo '# Jitsi Videobridge settings' > $CONFIG
            echo "JVB_HOSTNAME=$JVB_HOSTNAME" >> $CONFIG
            echo "JVB_PORT=$JVB_PORT" >> $CONFIG
            echo "JVB_SECRET=$JVB_SECRET" >> $CONFIG
            echo "JVB_OPTS=\"\"" >> $CONFIG

            # and we're done with debconf
            db_stop
        fi

        # we don't want to start the daemon as root
        if ! getent passwd jvb > /dev/null ; then
            useradd -r --user-group --shell /bin/bash --create-home -d /usr/share/jitsi-videobridge jvb
        fi

        # we create home folder only if it doesn't exist
        if [ ! -d /usr/share/jitsi-videobridge ]; then
            mkdir -p /usr/share/jitsi-videobridge
        fi

        # we claim the home folder of jvb in case it is owned by someone else
        OWNER=$(stat -c '%U' /usr/share/jitsi-videobridge)
        GROUP=$(stat -c '%G' /usr/share/jitsi-videobridge)
        if ! dpkg-statoverride --list /usr/share/jitsi-videobridge/* >/dev/null && [ "$OWNER:$GROUP" != "jvb:jvb" ]; then
            chown -R jvb:jvb /usr/share/jitsi-videobridge
            OWNER=jvb
            GROUP=jvb
        fi

        # die logz
        if [ ! -d /var/log/jitsi ]; then
            mkdir -p /var/log/jitsi
        fi
        chown -R jvb:jvb /var/log/jitsi

        # ensure videobridge is not running - it will be started at the end
        if [ "$(pidof java)" ]; then
            invoke-rc.d jitsi-videobridge stop || true
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
