# resources/install/debian/postrm
#!/bin/sh
# postrm script for jitsi-videobridge

set -e

# Load debconf
. /usr/share/debconf/confmodule


case "$1" in
    purge|remove)

        # Clear the debconf variable
        db_purge

        # stop the jitsi-videobridge
        if [ -x /etc/init.d/jitsi-videobridge ]; then
            invoke-rc.d jitsi-videobridge stop
        fi

        # clear user and group
        # ensure nothing is running as jvb user
        if ps -u jvb h; then
            `ps -u jvb -o pid h | xargs kill`
        fi
        # and then remove the user
        if getent passwd jvb > /dev/null ; then
            deluser jvb
        fi

        # clear logs
        if [ -d "/var/log/jitsi" ]; then
            rm -rf "/var/log/jitsi"
        fi
    ;;

    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
    ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

db_stop

exit 0
